# hosts/gcg-ax102/etc/systemd/system/openclaw-wonhoo.service.d/20-claude-cli-home.conf
[Service]
Environment=HOME=/opt/gcg/claude-profiles/sub-sonja
Environment=IS_SANDBOX=1
