# hosts/gcg-ax102/etc/systemd/system/openclaw-leon.service.d/20-claude-cli-home.conf
[Service]
Environment=HOME=/opt/gcg/claude-profiles/sub-jacques
Environment=IS_SANDBOX=1
