# hosts/gcg-ax102/etc/systemd/system/openclaw-varys.service.d/memory.conf
[Service]
MemoryMax=4G
MemoryHigh=3500M
