# hosts/gcg-ax102/etc/systemd/system/openclaw-daen.service.d/10-fleet-sender.conf
[Service]
Environment=FLEET_SENDER=daen
