# hosts/gcg-ax102/etc/systemd/system/openclaw-daen.service.d/critical6-lock.conf
[Unit]
RefuseManualStop=yes
RefuseManualStart=no
