# hosts/gcg-ax102/etc/systemd/system/openclaw-talos.service.d/20-claude-cli-home.conf
[Service]
Environment=HOME=/opt/gcg/claude-profiles/sub-sergei
Environment=IS_SANDBOX=1
