# hosts/gcg-ax102/etc/systemd/system/openclaw-talos.service.d/memory-node.conf
[Service]
# Override main unit's Environment=NODE_OPTIONS=--max-old-space-size=2048
# Raised to match 8G MemoryMax for concurrency=4 subagent workloads
Environment=NODE_OPTIONS=--max-old-space-size=4096
