# hosts/gcg-ax102/etc/systemd/system/openclaw-hector.service.d/20-env-secrets.conf
[Service]
Environment=GATEWAY_TOKEN=366fdc67805b524615c413ef891b0a894849192839caf027
Environment=HOOKS_TOKEN=9194e67a28a81ac548c95c396f368d4ed8e2ddb4c62573a58d73dadbccb036c3