# hosts/gcg-ax102/etc/systemd/system/openclaw-phil.service.d/20-env-secrets.conf
[Service]
Environment=GATEWAY_TOKEN=806b1ab79a1ecc851551c87982e490bd9ed2064ffdf60f6c
Environment=HOOKS_TOKEN=2d64317ac72d3987654ae99304e528c46a27a0b7cea68088e2fde0156325e3c9
