# hosts/gcg-ax102/etc/systemd/system/openclaw-chiron.service.d/10-secure-creds.conf
[Service]
ExecStartPre=
ExecStartPre=/usr/bin/python3 /opt/gcg/shared/scripts/validate-agent-config.py /opt/gcg/openclaw-chiron/openclaw.json
ExecStartPre=/opt/gcg/shared/bin/gcg-secret-bootstrap.sh chiron
LoadCredentialEncrypted=op-sa-token:/etc/credstore.encrypted/op-sa-token
LoadCredentialEncrypted=db-env:/etc/credstore.encrypted/db-env
LoadCredentialEncrypted=gateway-env:/etc/credstore.encrypted/gateway-env
ExecStart=
ExecStart=/bin/bash -c 'set -a; source /run/openclaw-chiron/env; exec /usr/bin/openclaw gateway'
