# hosts/gcg-ax102/etc/systemd/system/openclaw-yuri.service.d/20-claude-cli-home.conf
[Service]
Environment=HOME=/opt/gcg/claude-profiles/sub-vlada
Environment=IS_SANDBOX=1
