# hosts/gcg-ax102/etc/systemd/system/openclaw-socrates.service.d/10-secure-creds.conf
[Service]
Environment=FLEET_AGENT_NAME=socrates
Environment=HOOK_AGENT_ID=socrates
Environment=OPENCLAW_CONFIG_PATH=/opt/gcg/openclaw-socrates/openclaw.json
Environment=OPENCLAW_STATE_DIR=/opt/gcg/openclaw-socrates/state

LoadCredentialEncrypted=op-sa-token:/etc/credstore.encrypted/op-sa-token
LoadCredentialEncrypted=db-env:/etc/credstore.encrypted/db-env
LoadCredentialEncrypted=gateway-env:/etc/credstore.encrypted/gateway-env

RuntimeDirectory=openclaw-socrates
RuntimeDirectoryMode=0700

ExecStartPre=
ExecStartPre=/usr/bin/python3 /opt/gcg/shared/scripts/validate-agent-config.py /opt/gcg/openclaw-socrates/openclaw.json
ExecStartPre=/opt/gcg/shared/bin/gcg-secret-bootstrap.sh socrates
ExecStart=
ExecStart=/bin/bash -c 'set -a; source /run/openclaw-socrates/env; exec /usr/bin/openclaw gateway'
