# hosts/gcg-ax102/etc/systemd/system/gcg-daily-security-review.service.d/10-db-env.conf
[Service]
ExecStartPre=/bin/bash -c "systemd-creds decrypt /etc/credstore.encrypted/db-env /run/gcg-sec-review-env && chmod 600 /run/gcg-sec-review-env"
EnvironmentFile=-/run/gcg-sec-review-env
