# hosts/gcg-ax102/etc/systemd/system/gcg-inbox-poller@.service.d/10-backoff.conf
# Poke-storm fix (2026-06-13): re-poke backoff must exceed a typical turn (opus 30-90s)
# so the poller never spawns a 2nd concurrent turn before the first can claim/finish.
[Service]
Environment="INBOX_RETRY_BACKOFF=120 seconds"
