# hosts/gcg-ax102/etc/systemd/system/gcg-inbox-poller@.service.d/20-cascade-resume.conf
# W1 — cascade resume: allow the poller's Pass 4 hook execution to write to
# cascade state files and logs. task_cascade.py needs these paths at runtime.
#
# ExecStartPre ensures ReadWritePaths files exist — systemd mount namespacing
# fails with 226/NAMESPACE when a ReadWritePaths target doesn't exist on disk
# (incident 2026-06-25: gate_162547 — daen/nik/talos/viktor/vulcan all down).
[Service]
ExecStartPre=/bin/bash -c 'touch -a /opt/gcg/shared/cascade_task_state.json'
ReadWritePaths=/opt/gcg/shared/cascade_task_states
ReadWritePaths=/opt/gcg/shared/cascade_active_plans.json
ReadWritePaths=/opt/gcg/shared/cascade_task_state.json
ReadWritePaths=/opt/gcg/shared/cascade_verify_snapshots
ReadWritePaths=/opt/gcg/shared/cascade_agent_incidents.json
ReadWritePaths=/var/log
# Remove private tmp — hooks need to exchange state via shared filesystem paths
PrivateTmp=false
