# hosts/gcg-ax102/etc/systemd/system/gcg-inbox-poller@.service.d/30-liveness-aware.conf
# Fleet-wide liveness-aware relay (spine-v2 Phase 3.2, 2026-06-22).
# Absorbs wedge-recovery duty from gcg-inbox-stale-check into Pass 8 of gcg-inbox-poll.
# Canary (gcg-inbox-poller@daen) ran cleanly before this fleet-wide promotion.
# Phase 3.3: collapse gcg-inbox-stale-check once fleet soak passes.
[Service]
Environment=LIVENESS_AWARE=1
# Required for Pass 8's `systemctl restart openclaw-<agent>` under ProtectSystem=strict
ReadWritePaths=/run/systemd
