# hosts/gcg-ax102/etc/systemd/system/gcg-infra-snapshot.service
[Unit]
Description=GCG infra snapshot to git
After=network-online.target
[Service]
Type=oneshot
ExecStart=/opt/gcg/shared/bin/gcg-infra-snapshot.sh
TimeoutStartSec=15min
Nice=10
