# hosts/gcg-ax102/etc/systemd/system/gcg-infra-snapshot.timer
[Unit]
Description=Daily infra snapshot at 04:00 Dubai
[Timer]
OnCalendar=*-*-* 04:00:00 Asia/Dubai
AccuracySec=2min
Persistent=true
[Install]
WantedBy=timers.target
