# hosts/gcg-ax102/etc/systemd/system/gcg-verify-chain.service
[Unit]
Description=GCG Audit Log Chain Verifier
After=network.target

[Service]
Type=oneshot
ExecStart=/opt/gcg/shared/bin/verify-chain
User=root
StandardOutput=journal
StandardError=journal
