# hosts/gcg-ax102/etc/systemd/system/gcg-borg-backup.service
[Unit]
Description=GCG borg backup
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStartPre=-/opt/gcg/shared/bin/gcg-ax42-config-sync.sh
ExecStart=/opt/gcg/shared/bin/gcg-borg-backup.sh
LoadCredentialEncrypted=borg-passphrase:/etc/credstore.encrypted/borg-passphrase
LoadCredentialEncrypted=borg-key:/etc/credstore.encrypted/borg-key
Nice=10
IOSchedulingClass=idle
TimeoutStartSec=2h
