# hosts/gcg-ax102/etc/systemd/system/gcg-inbox-stale-check.timer
[Unit]
Description=Run GCG inbox stale-delivery monitor every 5 minutes

[Timer]
OnBootSec=3min
OnUnitActiveSec=5min
Unit=gcg-inbox-stale-check.service
Persistent=true

[Install]
WantedBy=timers.target
