# hosts/gcg-ax102/etc/systemd/system/gcg-inbox-stale-check.service
[Unit]
Description=GCG inbox stale-delivery monitor (alerts daen on rows unpicked >threshold)
After=network-online.target

[Service]
Type=oneshot
Environment=THRESHOLD=300
ExecStart=/opt/gcg/shared/bin/gcg-inbox-stale-check
# exit 1 = stale found + alerted (a normal outcome, not a unit failure)
SuccessExitStatus=0 1
Nice=10
MemoryMax=64M
NoNewPrivileges=true
# ProtectSystem=full (not strict): reconcile must re-materialize inbox files
# under /opt/gcg/openclaw-*/workspace/inbox. full keeps /usr /etc /boot read-only.
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
StateDirectory=gcg-inbox-stale-check
ReadWritePaths=/var/lib/gcg-inbox-stale-check /opt/gcg
