# hosts/gcg-ax102/etc/systemd/system/gcg-borg-backup.timer
[Unit]
Description=Run GCG borg backup nightly

[Timer]
OnCalendar=*-*-* 02:30:00 Asia/Dubai
RandomizedDelaySec=20m
Persistent=true
Unit=gcg-borg-backup.service

[Install]
WantedBy=timers.target
