# hosts/gcg-ax102/etc/systemd/system/gcg-inbox-janitor.timer
[Unit]
Description=Hourly GCG inbox janitor

[Timer]
OnBootSec=5min
OnUnitActiveSec=1h
Unit=gcg-inbox-janitor.service
Persistent=true

[Install]
WantedBy=timers.target
