# hosts/gcg-ax102/etc/systemd/system/gcg-secrets-pull.service
[Unit]
Description=GCG Secrets Pull from 1Password
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/etc/gcg/secrets-pull.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
