# hosts/gcg-ax102/etc/systemd/system/gcg-session-end-hook.service
[Unit]
Description=GCG Session End Hook — Extracts lessons from agent conversations
After=network.target

[Service]
Type=simple
User=root
EnvironmentFile=/opt/gcg/shared/credentials/session_end_hook.env
ExecStart=/usr/local/bin/gcg-session-end-hook-run
Restart=always
RestartSec=30s
StandardOutput=append:/opt/gcg/shared/logs/session_end_hook.log
StandardError=append:/opt/gcg/shared/logs/session_end_hook.log

[Install]
WantedBy=multi-user.target
