# hosts/gcg-ax102/etc/systemd/system/gcg-deep-audit.service
[Unit]
Description=GCG deep_audit_v2 sweep

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/gcg-deep-audit-wrapper.sh
Nice=10
IOSchedulingClass=idle
TimeoutStartSec=30m
