# hosts/gcg-ax102/etc/systemd/system/openclaw-daen.service
[Unit]
StartLimitBurst=5
StartLimitIntervalSec=300
Description=OpenClaw Agent - Daen HQ - The Architect
After=network.target docker.service

[Service]
EnvironmentFile=-/opt/gcg/shared/.env-db-secrets
EnvironmentFile=-/opt/gcg/shared/.env.global
Environment=FLEET_AGENT_NAME=daen
Type=simple
User=root
WorkingDirectory=/opt/gcg/openclaw-daen
EnvironmentFile=/opt/gcg/openclaw-daen/.env
EnvironmentFile=-/run/gcg/daen-secrets.env
EnvironmentFile=-/opt/gcg/shared/credentials/gateway.env
Environment=OPENCLAW_CONFIG_PATH=/opt/gcg/openclaw-daen/openclaw.json
Environment=OPENCLAW_STATE_DIR=/opt/gcg/openclaw-daen/state
Environment=HOOK_AGENT_ID=daen
Environment=NODE_OPTIONS=--max-old-space-size=8192
ExecStartPre=/bin/sh -c 'mkdir -p /run/gcg && OP_SERVICE_ACCOUNT_TOKEN=$(cat /opt/gcg/shared/credentials/op_service_account_token.txt) /usr/bin/op read "op://GCG Agent Fleet/Anthropic Daen/credential" | awk \'{ print "ANTHROPIC_API_KEY=" $0 }\' > /run/gcg/daen-secrets.env && chmod 600 /run/gcg/daen-secrets.env'
ExecStartPre=/bin/sh -c 'rm -f /opt/gcg/openclaw-daen/state/.openclaw/agents/*/sessions/*.jsonl.lock'
ExecStartPre=/usr/bin/python3 /opt/gcg/shared/scripts/validate-agent-config.py /opt/gcg/openclaw-daen/openclaw.json
ExecStart=/usr/bin/openclaw gateway
ExecStopPost=/opt/gcg/shared/scripts/session-pre-stop-flush.sh
MemoryMax=12G
MemoryHigh=10G
Restart=always
RestartSec=15
StandardOutput=append:/opt/gcg/openclaw-daen/logs/gateway.log
StandardError=append:/opt/gcg/openclaw-daen/logs/gateway.log

[Install]
WantedBy=multi-user.target
