# hosts/gcg-ax102/etc/systemd/system/openclaw-goku.service
[Unit]
StartLimitBurst=5
StartLimitIntervalSec=300
Description=OpenClaw Agent - Goku (Personal Health & Fitness Coach)
After=network.target docker.service

[Service]
EnvironmentFile=-/opt/gcg/shared/.env-db-secrets
EnvironmentFile=-/opt/gcg/shared/.env.global
Environment=FLEET_AGENT_NAME=goku
Environment=HOOK_AGENT_ID=goku
Type=simple
User=root
WorkingDirectory=/opt/gcg/openclaw-goku
EnvironmentFile=-/opt/gcg/shared/credentials/gateway.env
EnvironmentFile=/opt/gcg/openclaw-goku/.env
Environment=OPENCLAW_CONFIG_PATH=/opt/gcg/openclaw-goku/openclaw.json
Environment=OPENCLAW_STATE_DIR=/opt/gcg/openclaw-goku/state
Environment=NODE_OPTIONS=--max-old-space-size=2048
ExecStartPre=/usr/bin/python3 /opt/gcg/shared/scripts/validate-agent-config.py /opt/gcg/openclaw-goku/openclaw.json
ExecStart=/usr/bin/openclaw gateway
ExecStopPost=/opt/gcg/shared/scripts/session-pre-stop-flush.sh
MemoryMax=3G
MemoryHigh=2G
Restart=always
RestartSec=15
StandardOutput=append:/opt/gcg/openclaw-goku/logs/gateway.log
StandardError=append:/opt/gcg/openclaw-goku/logs/gateway.log

[Install]
WantedBy=multi-user.target
