# hosts/gcg-ax102/etc/systemd/system/openclaw-max.service
[Unit]
Description=OpenClaw Agent - Max (GCG - Olga Nechaieva)
After=network.target docker.service
StartLimitBurst=5
StartLimitIntervalSec=300

[Service]
EnvironmentFile=-/opt/gcg/shared/.env-db-secrets
EnvironmentFile=-/opt/gcg/shared/.env.global
Environment=FLEET_AGENT_NAME=max
Environment=HOOK_AGENT_ID=max
Type=simple
User=root
WorkingDirectory=/opt/gcg/openclaw-max
EnvironmentFile=/opt/gcg/openclaw-max/.env
Environment=OPENCLAW_CONFIG_PATH=/opt/gcg/openclaw-max/openclaw.json
Environment=OPENCLAW_STATE_DIR=/opt/gcg/openclaw-max/state
Environment=NODE_OPTIONS=--max-old-space-size=2048
ExecStartPre=/usr/bin/python3 /opt/gcg/shared/scripts/validate-agent-config.py /opt/gcg/openclaw-max/openclaw.json
ExecStart=/usr/bin/openclaw gateway
ExecStopPost=/opt/gcg/shared/scripts/session-pre-stop-flush.sh
MemoryMax=3G
MemoryHigh=3G
Restart=always
RestartSec=10
StandardOutput=append:/opt/gcg/openclaw-max/logs/gateway.log
StandardError=append:/opt/gcg/openclaw-max/logs/gateway.log

[Install]
WantedBy=multi-user.target
