# hosts/gcg-ax102/etc/systemd/system/openclaw-niccolo.service.d/20-claude-cli-home.conf
[Service]
Environment=HOME=/opt/gcg/claude-profiles/sub-ishan
Environment=IS_SANDBOX=1
