# hosts/gcg-ax102/etc/systemd/system/openclaw-niccolo.service.d/memory-node.conf
[Service]
# Override main unit's 2048 → 4096 (paired with 8G MemoryMax for fleet-wide analysis workloads)
Environment=NODE_OPTIONS=--max-old-space-size=4096
