# hosts/gcg-ax102/etc/systemd/system/openclaw-tom.service.d/20-env-secrets.conf
[Service]
Environment=GATEWAY_TOKEN=5db1e9e2d761a934d916e8975f3cb23cf074278a7900e58f
Environment=HOOKS_TOKEN=4ece75014a72ac5578d02e6e500bb2cbb1825af440f9aca2e7ed9c7766783ec1