# hosts/gcg-ax102/etc/systemd/system/openclaw-malik.service.d/10-fleet-sender.conf
[Service]
Environment=FLEET_SENDER=malik
