# hosts/gcg-ax102/etc/systemd/system/openclaw-angela.service.d/20-claude-cli-home.conf
[Service]
Environment=HOME=/opt/gcg/claude-profiles/sub-petr
Environment=IS_SANDBOX=1
